# deploy_configs/supervisor/development.ini
[program:{{ project_name }}-development]
directory = <path to {{ project_name }} root directory>
command = <path to virtualenv uwsgi binary> --ini <path to {{project_name}} root directory>/deploy_configs/uwsgi/development.ini
autostart = true
autorestart = true

